# deploy/basic_flask_template.ini
[program:basic_flask_template]
command=/bin/bash -c "cd ~/basic_flask_template && ./run_uwsgi.sh"
autostart=yes
autorestart=yes
stopsignal=INT
stopasgroup=true
killasgroup=true

