# deploy/basic_flask_template.ini
[program:basic_flask_template]
chdir=/home/bewi/basic_flask_template
home=/home/bewi/basic_flask_template
pythonpath=/home/bewi/basic_flask_template/ENV/bin/python
virtualenv=/home/bewi/basic_flask_template/ENV+
command=/bin/bash -c '~/basic_flask_template/ENV/bin/uwsgi uwsgi.ini'
directory=/home/bewi/basic_flask_template
autostart=yes
autorestart=yes
stopsignal=INT
stopasgroup=true
killasgroup=true

