# deploy/helloworld.ini
[program:project_name-uwsgi]
command=/bin/bash -c "cd ~/project_name && ./run_uwsgi.sh"
autostart=yes
autorestart=yes

